# pidoors/pidoors.service
[Unit]
Description=PiDoors Access Control Service
After=network.target mysql.service
Wants=network.target

[Service]
Type=simple
User=pidoors
Group=pidoors
WorkingDirectory=/opt/pidoors
ExecStart=/opt/pidoors/venv/bin/python3 /opt/pidoors/pidoors.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/pidoors/cache
ReadOnlyPaths=/opt/pidoors/conf /opt/pidoors/pidoors.py

# Allow GPIO access
SupplementaryGroups=gpio

[Install]
WantedBy=multi-user.target
