# pidoors/pidoors.service
[Unit]
Description=PiDoors Access Control Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pidoors
Group=pidoors
RuntimeDirectory=pidoors
WorkingDirectory=/run/pidoors
ExecStart=/opt/pidoors/venv/bin/python3 /opt/pidoors/pidoors.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening (NoNewPrivileges removed to allow sudo for updates)
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/pidoors/cache /opt/pidoors/conf

# Allow GPIO access
SupplementaryGroups=gpio

[Install]
WantedBy=multi-user.target
